# deploy/sdrtrunk.service
# SDRTrunk Kennebec systemd service - automatic restart on crash for unattended operation.
#
# Install:
#   1. Adjust ExecStart path (and User) below for your installation.
#   2. sudo cp deploy/sdrtrunk.service /etc/systemd/system/
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable --now sdrtrunk
#
# Health monitoring: the application serves http://127.0.0.1:8080/health (port configurable via
# SDRTRUNK_API_PORT) returning HTTP 200 when healthy and 503 when a tuner or stream is degraded.
# POST http://127.0.0.1:8080/restart exits the application so systemd restarts it; set
# SDRTRUNK_API_TOKEN to require an Authorization: Bearer token for the restart endpoint.

[Unit]
Description=SDRTrunk Kennebec radio monitoring
After=network.target

[Service]
Type=simple
# Run as a dedicated user with access to USB devices (plugdev group on Debian/Ubuntu)
User=sdrtrunk
# Headless mode - no GUI required
ExecStart=/opt/sdrtrunk/bin/sdrtrunk --headless
Restart=always
RestartSec=10
# Prevent rapid crash-loops from spinning: at most 5 restarts in 10 minutes
StartLimitIntervalSec=600
StartLimitBurst=5
# Optional REST API configuration
#Environment=SDRTRUNK_API_PORT=8080
#Environment=SDRTRUNK_API_TOKEN=change-me

[Install]
WantedBy=multi-user.target
